# api/supervisor/supervisord.conf
[supervisord]
logfile=/home/mainsail/klipper_logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10 ; (num of main logfile rotation backups;default 10)
loglevel=info ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true ; (start in foreground if true;default false)
minfds=1024 ; (min. avail startup file descriptors;default 1024)
minprocs=200 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock

[program:klipper]
command=/home/mainsail/klippy-env/bin/python /home/mainsail/klipper/klippy/klippy.py /home/mainsail/klipper_config/printer.cfg -l /home/mainsail/klipper_logs/klippy.log -a /tmp/klippy_uds
user=mainsail
process_name=klipper
directory=/home/mainsail
autostart=true
autorestart=true
redirect_stderr=true

[program:moonraker]
command=/home/mainsail/moonraker-env/bin/python /home/mainsail/moonraker/moonraker/moonraker.py -l /home/mainsail/klipper_logs/moonraker.log -c /home/mainsail/klipper_config/moonraker.conf
user=mainsail
process_name=moonraker
directory=/home/mainsail
autostart=true
autorestart=true
redirect_stderr=true
